# deploy/systemd/maillogsentinel.timer
[Unit]
Description=Run MailLogSentinel hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
